# dedup.cron
#!/bin/sh

bedup dedup --defrag /mnt/HDD
bedup dedup /mnt/SSD

