# dedup.cron
#!/bin/sh

bedup dedup --defrag /mnt/HDD > /dev/null
bedup dedup /mnt/SSD > /dev/null
