# defrag.cron
#!/bin/sh
btrfs_root_volumes=/mnt/HDD
for btrfs_root_volume in $btrfs_root_volumes; do
  volumes=`btrfs subvolume list "$btrfs_root_volume" | awk '{print $9}'`
  for vol in $volumes; do
    find "$btrfs_root_volume"/$vol -xdev \( -type d -o -type f \) -print -exec btrfs file defrag "{}" \;
  done
done
