# fstrim.cron
#!/bin/sh

for mountpoint in / /boot /home /var/lib/mysql ; do
  /usr/sbin/fstrim "$mountpoint"
done
